# astrodb-explorer/docs/nina-ingest.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job for the NINA ingest worker on office-mac.

  Install:
    cp docs/nina-ingest.launchd.plist ~/Library/LaunchAgents/com.briancarter.nina-ingest.plist
    launchctl load   ~/Library/LaunchAgents/com.briancarter.nina-ingest.plist   # start + run at login
    launchctl unload ~/Library/LaunchAgents/com.briancarter.nina-ingest.plist   # stop

  Logs: ~/Library/Logs/nina-ingest.log  (stdout + stderr)
  Reads PG creds from astrodb-explorer/.env.local; NINA_HOST defaults to 100.65.96.18.

  Absolute node path (nvm v22.19.0) is used so launchd needs no shell init.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briancarter.nina-ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/briancarter/.nvm/versions/node/v22.19.0/bin/node</string>
        <string>--env-file=.env.local</string>
        <string>--import</string>
        <string>tsx</string>
        <string>scripts/nina-ingest.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/Office-SSD/Astronomy/astro-tools/astrodb-explorer</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart the worker if it ever exits; it self-heals NINA disconnects
         internally, so an exit means something unexpected. -->
    <key>KeepAlive</key>
    <true/>

    <!-- Don't hammer restart loops if it crashes immediately at boot. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/briancarter/Library/Logs/nina-ingest.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/briancarter/Library/Logs/nina-ingest.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
